# slangbrain.socket
[Socket]
ListenStream=80
ListenStream=443
KeepAlive=true
KeepAliveTimeSec=180

[Install]
WantedBy=sockets.target
